# overlay/etc/s6-overlay/s6-rc.d/vsftpd/finish
#!/command/with-contenv sh

set -euo pipefail
export BF_E="${PWD##*/}/$(basename ${0})"


#======================================================================================================================
# Show helpful log message.
#======================================================================================================================

bf-svc-finish -t
