# overlay/etc/s6-overlay/s6-rc.d/vsftpd/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

nu -c "use bf-ftps run ; run preflight"


#======================================================================================================================
# Run vsftpd server.
#======================================================================================================================

vsftpd ${BF_FTPS_VSFTPD_CONF}
