# overlay/etc/s6-overlay/s6-rc.d/vsftpd/run
#!/usr/bin/with-contenv sh

set -euo pipefail
export BF_E="${PWD##*/}/$(basename ${0})"


#======================================================================================================================
# Run vsftpd server.
#======================================================================================================================

[[ ! -f ${FTPS_VSFTPD_CONF} ]] && bf-error "vsftpd configuration file not found."

bf-echo "Starting vsftpd."
vsftpd ${FTPS_VSFTPD_CONF}
