# overlay/etc/s6-overlay/s6-rc.d/vsftpd-log/run
#!/usr/bin/nu

use bf
use bf-s6
bf env load -x vsftpd-log

# Forward vsftpd log to Docker
def main [...args] { bf-s6 log forward vsftpd-log (bf env FTPS_VSFTPD_LOG_FILE) }
