# overlay/etc/services.d/vsftpd-log/finish
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Show helpful log message.
#======================================================================================================================

bf-svc-finish -s ${PWD##*/}
