# overlay/etc/services.d/vsftpd-log/run
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Forward VSFTPD error log to Docker
#======================================================================================================================

bcg-forward "vsftpd-log" "/var/log/vsftpd.log"
