# overlay/etc/services.d/vsftpd-log/run
#!/usr/bin/with-contenv sh

set -euo pipefail
export BF_E="${PWD##*/}/$(basename ${0})"


#======================================================================================================================
# Forward vsftpd error log to Docker.
#======================================================================================================================

bf-forward ${PWD##*/} ${FTPS_VSFTPD_LOG_FILE}
