# overlay/etc/services.d/vsftpd/finish
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Use base executable to terminate all services
#======================================================================================================================

_terminate
