# overlay/etc/services.d/vsftpd/run
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Run VSFTPD server
#======================================================================================================================

_echo "Starting VSFTPD" 2>&1 | tee /var/log/vsftpd.log
/usr/sbin/vsftpd /etc/vsftpd/vsftpd.conf
