# overlay/etc/services.d/vsftpd/run
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Run VSFTPD server
#======================================================================================================================

CONF=/etc/vsftpd/vsftpd.conf
[[ ! -f ${CONF} ]] && _error "VSFTPD configuration file not found."

_echo "Starting VSFTPD"
vsftpd ${CONF}
