# overlay/etc/services.d/vsftpd/run
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Run VSFTPD server
#======================================================================================================================

CONF=/etc/vsftpd/vsftpd.conf
[[ ! -f ${CONF} ]] && bcg-error "vsftpd configuration file not found."

bcg-echo "Starting vsftpd"
vsftpd ${CONF}
