# overlay/etc/services.d/vsftpd/run
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Run VSFTPD server
#======================================================================================================================

_echo "Starting VSFTPD"
exec /usr/sbin/vsftpd /etc/vsftpd/vsftpd.conf
